# vol-osd.service
[Unit]
Description=vol-osd — per-app volume OSD daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
Environment=LD_PRELOAD=/usr/lib/x86_64-linux-gnu/libgtk4-layer-shell.so.0
ExecStart=%h/.local/bin/vol-osd
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
